# .config/autostart/sxhkd.desktop
[Desktop Entry]
Name=sxhkd
TryExec=sxhkd
Exec=sxhkd
Type=Application
